# cmd/m17-gateway/debian/postinst
#!/bin/bash

set -e

# Source debconf library
. /usr/share/debconf/confmodule

# Create service user and group if they don't exist
if ! getent group m17-gateway >/dev/null; then
    groupadd --system m17-gateway
fi

if ! getent passwd m17-gateway >/dev/null; then
    useradd --system --gid m17-gateway --home-dir /opt/m17/m17-gateway --shell /bin/false m17-gateway
fi

# Create directories
mkdir -p /opt/m17/m17-gateway
mkdir -p /opt/m17/etc

# Set permissions
chown -R m17-gateway:m17-gateway /opt/m17
chmod 755 /opt/m17/m17-gateway
chmod 755 /opt/m17/etc

# Set config file permissions (readable by service user)
chmod 644 /opt/m17/etc/m17-gateway.ini

# Make binary executable
chmod +x /opt/m17/m17-gateway/m17-gateway

# Get the callsign from debconf
db_get m17-gateway/callsign || true
CALLSIGN="$RET"

# Update the configuration file with the user's callsign
if [ -n "$CALLSIGN" ]; then
    sed -i "s/CALLSIGN_PLACEHOLDER/$CALLSIGN/" /opt/m17/etc/m17-gateway.ini
    echo "Configuration updated with callsign: $CALLSIGN"
else
    echo "No callsign provided. Please edit /opt/m17/etc/m17-gateway.ini manually."
fi

# Reload systemd and enable the service
systemctl daemon-reload
systemctl enable m17-gateway.service

echo
echo "M17 Gateway installed successfully!"
echo "Configuration file: /opt/m17/etc/m17-gateway.ini"
echo "View logs with: journalctl -u m17-gateway"
echo "Start with: sudo systemctl start m17-gateway"
echo "Check status with: sudo systemctl status m17-gateway"

# Clean up debconf database
db_stop
