# cmd/m17-gateway/debian/postinst
#!/bin/bash

set -e

# Create service user and group if they don't exist
if ! getent group m17-gateway >/dev/null; then
    groupadd --system m17-gateway
fi

if ! getent passwd m17-gateway >/dev/null; then
    useradd --system --gid m17-gateway --groups dialout,gpio --home-dir /opt/m17/m17-gateway --shell /bin/false m17-gateway
fi

# Create directories
mkdir -p /opt/m17/m17-gateway

# Set permissions
chown m17-gateway:m17-gateway /opt/m17/m17-gateway
chmod 755 /opt/m17/m17-gateway

# Make binary executable
chmod +x /opt/m17/m17-gateway/m17-gateway

# Only modify m17-gateway.ini if it doesn't exist, i.e. this is an upgrade
if [ ! -f /etc/m17-gateway.ini ]; then
    # Prompt for callsign if we're in an interactive terminal
    CALLSIGN=""
    if [ -t 0 ]; then
        echo
        echo "M17 Gateway Configuration"
        echo "========================="
        echo "Please enter your amateur radio callsign (e.g., N1ADJ or N1ADJ G):"
        echo "(Leave blank if you want to configure this manually later)"
        read -r CALLSIGN
    fi

    # Update the configuration file with the user's callsign
    if [ -n "$CALLSIGN" ]; then
        # Use awk to replace the placeholder - more robust than sed for this case
        awk -v callsign="$CALLSIGN" '
            /^Callsign=CALLSIGN_PLACEHOLDER/ { print "Callsign=" callsign; next }
            { print }
        ' /etc/m17-gateway.ini.sample > /etc/m17-gateway.ini
        rm /etc/m17-gateway.ini.sample
        echo "Configuration updated with callsign: $CALLSIGN"
    else
        mv /etc/m17-gateway.ini.sample /etc/m17-gateway.ini
        echo "No callsign provided. Please edit /etc/m17-gateway.ini manually."
        echo "Change 'CALLSIGN_PLACEHOLDER' to your actual callsign."
    fi
fi

# Set config file ownership and permissions (readable by service user)
chown m17-gateway:m17-gateway /etc/m17-gateway.ini
chmod 644 /etc/m17-gateway.ini

# Reload systemd and enable the service
systemctl daemon-reload
systemctl enable m17-gateway.service

echo
echo "M17 Gateway installed successfully!"
echo "Configuration file: /etc/m17-gateway.ini"
echo "View logs with: journalctl -u m17-gateway"
echo "Start with: sudo systemctl start m17-gateway"
echo "Check status with: sudo systemctl status m17-gateway"
